# installation/config/systemd/snort.service
# ============================================================================
# SIEM AFRICA - Service systemd Snort
# ============================================================================
# Fichier : /etc/systemd/system/snort.service
#
# HARD-WON FIXES V1 :
#   - PIDFile dans /var/log/ (persistant, pas /var/run/ éphémère)
#   - Pas de ProtectSystem=strict (cause CHDIR error 200)
#   - Pas de PrivateTmp=yes (bloque accès logs)
# ============================================================================

[Unit]
Description=SIEM Africa - Snort NIDS
Documentation=https://www.snort.org/documents
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=forking
User=snort
Group=snort

# PID file persistant (hard-won fix V1)
PIDFile=/var/log/snort/snort.pid

# Variables d'env (INTERFACE est défini dans /etc/default/snort)
EnvironmentFile=-/etc/default/snort

# Commande de démarrage
ExecStart=/usr/sbin/snort \
    -D \
    -u snort \
    -g snort \
    -c /etc/snort/snort.conf \
    -i ${INTERFACE} \
    -l /var/log/snort \
    --pid-path /var/log/snort \
    --create-pidfile

# Reload à chaud des règles
ExecReload=/bin/kill -HUP $MAINPID

# Arrêt propre
ExecStop=/bin/kill -TERM $MAINPID
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart automatique
Restart=on-failure
RestartSec=10s

# Durcissement (SANS ProtectSystem=strict ni PrivateTmp=yes)
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Capabilities réseau nécessaires
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SETUID CAP_SETGID CAP_DAC_OVERRIDE CAP_SYS_CHROOT
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
